# etc/supervisord.d/apache.ini
[program:apache-httpd]
command=/usr/sbin/httpd -DFOREGROUND
autostart=true
autorestart=unexpected
exitcodes=0
