# etc/supervisord.d/php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm -F
autostart=true
autorestart=unexpected
exitcodes=0
